# etc/dotfiles/dot-config/systemd/user/wss-bridge-chrome.service
[Unit]
Description=wss-bridge to mcp.danenberg.ai (chrome-cdp-mcp --stdio on laptop, Chrome :9223)
After=network-online.target
Wants=network-online.target

[Service]
Type=simple
ExecStart=%h/bin/wss-bridge-chrome-launcher
Restart=on-failure

[Install]
WantedBy=default.target
